# res/etc/profile.d/cstbox.sh
#!/bin/bash
export PATH="/opt/cstbox/bin:$PATH"
